# systemd/.config/systemd/user/update-drive.path
[Unit]
Description= Updates Google Drive upon change.

[Path]
PathModified=/home/erick/Google Drive/

[Install]
WantedBy=default.target
