# systemd/.config/systemd/user/update-drive.timer
[Unit]
Description= Updates Google Drive hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
