# systemd/.config/systemd/user/update-drive.service
[Unit]
Description= Updates Google Drive


[Service]
Type=oneshot
ExecStart=/home/erick/.local/bin/update_gdrive.sh
